# omnilink-bridge/rootfs/etc/services.d/omnilink-bridge/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on
# ==============================================================================

bashio::log.info "Starting OmniLink Bridge..."

CONFIG_PATH=/data/options.json

jq -r 'to_entries|map("\(.key)=\(.value|tostring)")|.[]' $CONFIG_PATH > /app/OmniLinkBridge.ini

ARGS=
if $(bashio::config.true 'support_send_logs_to_developer'); then
  ARGS=-ld
fi

mono /app/OmniLinkBridge.exe -i -c /app/OmniLinkBridge.ini -lf disable $ARGS
